# addons/point_of_sale/tools/posbox/overwrite_before_init/etc/ld.so.preload
# License LGPL-3.0 or later (http://www.gnu.org/licenses/lgpl).
#/usr/lib/arm-linux-gnueabihf/libarmmem.so
